# docs/systemd/api.service
[Unit]
Description=Rubika Bot SaaS API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rubika-bot-saas
EnvironmentFile=/opt/rubika-bot-saas/.env
ExecStart=/usr/bin/env bash -lc 'cd /opt/rubika-bot-saas && bash scripts/run_api.sh'
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
